# scripts/schoolapp-podman.service
# SchoolApp Podman Auto-Start Service
# This systemd service ensures Podman containers start automatically on boot
# Note: Podman is daemonless, so no podman.service dependency is needed

[Unit]
Description=SchoolApp Podman Containers
Documentation=https://github.com/your-repo/SchoolApp
# Podman is daemonless, so we don't require a podman.service
# But we want network to be available
Wants=network-online.target
After=network-online.target
# Optional: wait for user session if running rootless
After=user@%i.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/path/to/your/SchoolApp
ExecStart=/path/to/your/SchoolApp/scripts/podman-auto-start.sh start
ExecStop=/path/to/your/SchoolApp/scripts/podman-auto-start.sh stop
ExecReload=/path/to/your/SchoolApp/scripts/podman-auto-start.sh restart
TimeoutStartSec=300
TimeoutStopSec=60

# Restart policy
Restart=on-failure
RestartSec=10

# Security settings
# For rootless Podman, use your user account instead of root
User=%i
Group=%i
# For rootful Podman, uncomment below and comment above
# User=root
# Group=root

# Environment
Environment=COMPOSE_PROJECT_NAME=schoolapp
# Ensure Podman can access user's container storage (for rootless)
Environment=XDG_RUNTIME_DIR=/run/user/%i

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=schoolapp-podman

[Install]
WantedBy=multi-user.target

